# files/rootfs/etc/services.d/openbox/run
#!/usr/bin/execlineb -P

fdmove -c 2 1

foreground { s6-svwait -u /var/run/s6/services/crashplan-server }
foreground { s6-svwait -u /var/run/s6/services/tigervnc }
foreground { s6-ftrig-wait /var/run/s6/services/tigervnc/event S }

export DISPLAY :0
export HOME /root

/usr/bin/openbox-session
