# files/rootfs/etc/services.d/novnc/run
#!/usr/bin/execlineb -P

fdmove -c 2 1

/opt/novnc/utils/launch.sh --listen 6080 --vnc localhost:5900
